# mods/sshd/root/etc/s6-overlay/s6-rc.d/svc-mod-dropbear/run
#!/usr/bin/with-contenv bash
exec /usr/sbin/dropbear -F -E -p 22 \
    -r /config/.ssh/dropbear_ecdsa_host_key \
    -r /config/.ssh/dropbear_ed25519_host_key \
    -s \
    -e
