# mods/sshd/root/etc/s6-overlay/s6-rc.d/init-mod-dropbear/run
#!/usr/bin/with-contenv bash

# Generate persistent host keys in ~/.ssh (survives container restarts)
mkdir -p /config/.ssh
if [ ! -f /config/.ssh/dropbear_ecdsa_host_key ]; then
    dropbearkey -t ecdsa -f /config/.ssh/dropbear_ecdsa_host_key
fi
if [ ! -f /config/.ssh/dropbear_ed25519_host_key ]; then
    dropbearkey -t ed25519 -f /config/.ssh/dropbear_ed25519_host_key
fi

chmod 700 /config/.ssh
touch /config/.ssh/authorized_keys
chmod 600 /config/.ssh/authorized_keys
chown -R abc:abc /config/.ssh

# Expose key container env vars to SSH login sessions.
# Dropbear starts a clean shell that doesn't inherit the container environment,
# so vars like TZ and XDG_RUNTIME_DIR would otherwise be missing, breaking
# timezone display and tools like zellij that locate sockets via XDG_RUNTIME_DIR.
mkdir -p /etc/profile.d
cat > /etc/profile.d/devtop-env.sh <<EOF
export TZ=${TZ:-UTC}
export XDG_RUNTIME_DIR=${XDG_RUNTIME_DIR:-/config/.XDG}
EOF
